# deploy/aicrash.service
[Unit]
Description=AICrash.news - AI Industry Crash News Monitor
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=aicrash
Group=aicrash
WorkingDirectory=/opt/aicrash
EnvironmentFile=/opt/aicrash/.env
ExecStart=/usr/bin/node server/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
